# macmini/com.cdw.refresh-tab.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cdw.refresh-tab</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/toast-it-devices/cdw-slack-bot/refresh_cdw_tab.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/cdw-refresh-tab.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/cdw-refresh-tab-error.log</string>
</dict>
</plist>
